# ola/rootfs/etc/s6-overlay/s6-rc.d/ola/run
#!/command/with-contenv sh
exec s6-setuidgid ola /usr/bin/olad --config-dir /etc/ola -f
