# ola/rootfs/etc/s6-overlay/s6-rc.d/ola/run
#!/command/with-contenv sh
chown -R ola:ola /config
exec s6-setuidgid ola /usr/bin/olad --config-dir /config -f
